# src/lib/config/test/program_invalid_value_exitcode.ini
[program:a]
command = "/bin/ls"
exitcodes = "asdf"