# src/lib/config/test/program.ini
[program:a]
command = "/bin/ls"
numprocs = 3
environment = "A=1,B=2"
exitcodes = "1,2,3"
umask = "0222"
autostart = false
autorestart = never
stopsignal = KILL