# src/lib/config/test/program_invalid_value_environment.ini
[program:a]
command = "/bin/ls"
environment = "A=1,B"