# src/lib/config/test/program_invalid_value_umask.ini
[program:a]
command = "/bin/ls"
umask = "asdf"