# src/lib/config/test/program_invalid_value_bool.ini
[program:a]
command = "/bin/ls"
autostart = "ff"