# src/lib/config/test/program_invalid_value_u32.ini
[program:a]
command = "/bin/ls"
numprocs = "three"