# src/lib/config/test/program_invalid_value_stopsignal.ini
[program:a]
command = "/bin/ls"
stopsignal = asdf