# src/lib/config/test/program_invalid_value_autorestart.ini
[program:a]
command = "/bin/ls"
autorestart = "asdf"