# test/conf_user.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:conf_user]
command = "/usr/bin/id -u"
stdout_logfile = "/tmp/conf_user.log"
user = "test"