# test/conf_environment.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:conf_environment]
command = "/usr/bin/env"
stdout_logfile = "test/conf_environment.log"
environment = "A=1,B=2"