# test/conf_autorestart_never.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:conf_autorestart_never]
command = "/usr/bin/sleep 2"
autostart = true
autorestart = never
exitcodes = "1"

