# test/conf_autorestart_always.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:conf_autorestart_always]
command = "/usr/bin/sleep 2"
autostart = true
autorestart = always
exitcodes = "0,1"

