# test/conf_exitcodes.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:conf_exitcodes]
command = "/usr/bin/sleep 2"
autostart = true
autorestart = unexpected
exitcodes = "0"

