# test/conf_startsecs.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:conf_startsecs]
command = "/usr/bin/sleep 10"
autostart = true
startsecs = 5

