# test/conf_directory.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:conf_directory]
directory = "/tmp"
command = "/usr/bin/pwd"
stdout_logfile = "/tmp/conf_directory.log"