# test/conf_autostart.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:conf_autostart]
command = "/usr/bin/tail -f"
autostart = true
