# test/conf_numprocs.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:conf_numprocs]
command = "/usr/bin/sleep 10"
numprocs = 3