# test/conf_command.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:conf_command]
command = "/home/mypark/supervisor/test/conf_command.sh arg1 arg2 arg3 arg4"
stdout_logfile = "test/conf_command.log"