# test/conf_startsecs2.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:conf_startsecs2]
command = "/usr/bin/sleep 4"
autostart = true
startsecs = 5

