# test/conf_autorestart_unexpected.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:conf_autorestart_unexpected]
command = "/usr/bin/sleep 2"
autostart = true
autorestart = unexpected
exitcodes = "1"

