# test/restart_simple_success2.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:tailf]
command = "tail -f"
autostart = true