# test/start_simple_fail.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:dies]
command = "/usr/bin/ls /"
user = "test"