# test/status_before_begin2.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:a]
command = "/usr/bin/ls"

[program:b]
command = "/usr/bin/ls"

[program:c]
command = "/usr/bin/ls"

[program:d]
command = "/usr/bin/ls"