# test/stderr_logfile.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:eecho]
directory = "."
command = "/bin/bash test/eecho.sh Hello"
stderr_logfile = "/tmp/eecho.err"