# test/stop_signal1.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:sig_echo]
directory = "."
command = "/bin/bash test/echo_signal.sh"
stdout_logfile = "/tmp/sig_echo.log"
stopsignal = INT