# test/stopwaitsecs2.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:ign_term]
command = "/bin/bash test/ign_term.sh"
directory = "."
stopsignal = TERM
stopwaitsecs = 3