# test/umask.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:touch]
command = "touch /tmp/somefile.txt"
umask = 700