# test/update_simple_origin.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:a]
command = "/usr/bin/ls"