# test/update_simple_modified.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:b]
command = "/usr/bin/ls"