# general.ini
[general]
sockfile="/tmp/supervisor.sock"

[program:a]
command = "/usr/bin/ls"
numprocs = 3
environment = "A=1,B=2"
exitcodes = "1,2,3"
umask = "0222"
autostart = false
autorestart = never
stopsignal = KILL

[program:tail]
command = "/usr/bin/tail -f"
numprocs = 2
umask = "0222"
autostart = false
autorestart = never
stopsignal = TERM


[program:sleep]
command = "/usr/bin/sleep 10"
numprocs = 1
umask = "0222"
autorestart = never
stopsignal = TERM
