# general.ini
[general]
sockfile="/tmp/taskmaster.sock"

[program:a]
command = "/usr/bin/sleep 3"
startsecs = 5
startretries = 2

[program:ign_term]
command = "/root/supervisor/test/a.out"
stopwaitsecs = 5

[program:test_signal]
command = "/root/supervisor/test/b.out"
stopsignal = INT

[program:ls-err]
command = "ls /thereisnofile"
stderr_logfile = "/tmp/ls-err.log"

[program:ls-ok]
command = "ls /"
stdout_logfile = "/tmp/ls-ok.log"
